# .config/systemd/user/ollama.service
[Unit]
Description=Ollama Serve Daemon

[Service]
ExecStart=/usr/bin/ollama serve
Restart=on-failure

[Install]
WantedBy=default.target
